# rc.local
/bin/sh /etc/config/sta.sh & > /dev/null 2>&1
