# rc.local
/bin/sh /etc/config/owrf.sh & > /dev/null 2>&1
